# home/systemd/dot-config/systemd/user/hintsd.service
[Unit]
Description=Hints daemon
[Service]
Type=simple
ExecStart=/home/yatin/.local/bin/hintsd
Restart=always
[Install]
WantedBy=default.target
