# home/systemd/dot-config/systemd/user/hyprshade-grayscale.service
[Unit]
Description=Enable hyprshade grayscale

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprshade on grayscale
