# systemd/.config/systemd/user/email.timer
[Unit]
Description=Sync Email

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
