# systemd/.config/systemd/user/email.timer
[Unit]
Description=Sync Email

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
