# systemd/.config/systemd/user/email.service
[Unit]
Description=Sync Email

[Service]
Type=oneshot
ExecStart=/usr/bin/zsh /home/nitay/code/scripts/syncemail

[Install]
WantedBy=default.target
