# systemd/.config/systemd/user/newsboat.timer
[Unit]
Description=Update newsboat periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
